# face-assist/rootfs/etc/services.d/face-assist/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Face Assist service
# ==============================================================================

declare -a options

# Ensure directory permissions
bashio::log.info "Setting up directories..."
mkdir -p /share/face-assist/models \
         /share/face-assist/faces \
         /share/face-assist/results

# Set the log level
export LOG_LEVEL="$(bashio::config 'log_level')"

# Wait for Hailo device
bashio::log.info "Waiting for Hailo device..."
TIMEOUT=30
count=0
while [ ! -e /dev/hailo0 ]; do
    sleep 1
    count=$((count + 1))
    if [ $count -ge $TIMEOUT ]; then
        bashio::log.warning "Hailo device not found after ${TIMEOUT} seconds, continuing anyway..."
        break
    fi
done

cd /usr/src/app || bashio::exit.nok "Could not change to app directory"

# Start the Flask application
bashio::log.info "Starting Face Assist..."
exec python3 app.py
